# com.user.whisper-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.whisper-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/blake/projects/whisper-transcribe-ui/.venv/bin/python</string>
        <string>-m</string>
        <string>app.daemon.whisper_daemon</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/blake/projects/whisper-transcribe-ui</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/blake/projects/whisper-transcribe-ui/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/blake/projects/whisper-transcribe-ui</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/whisper-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/whisper-daemon.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
